# whatsapp-auth-restore.service
[Unit]
Description=WhatsApp Auth Session Restore Service
After=network.target
Before=nord-laundry-bot.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/nord-laundry-bot
ExecStart=/opt/nord-laundry-bot/auto-restore-auth.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

# Запускается только если папка авторизации отсутствует или повреждена
ConditionPathExists=!/opt/nord-laundry-bot/.wwebjs_auth
ConditionPathExists=!/opt/nord-laundry-bot/.wwebjs_auth/session-*

[Install]
WantedBy=multi-user.target
